# file_duplicates/rc.local
#!/bin/bash -e


bash  /home/mac/Desktop/exp/brightness/scripts
exit 0
